# deploy/macos/com.tradingbot.start.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 每日 08:30 啟動交易 bot。__PROJECT_DIR__ 由 install.sh 換成實際路徑。 -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingbot.start</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_DIR__/deploy/macos/start_bot.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- 每天 08:30（本機時間 = UTC+8，與 bot 內建 Asia/Taipei 排程一致） -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- 不在載入/登入時自動啟動，只在 08:30 觸發；崩潰重啟由 supervisor 處理，故不設 KeepAlive -->
    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PYTHONUTF8</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/launchd_start.out.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/launchd_start.err.log</string>
</dict>
</plist>
