# deploy/macos/com.tradingbot.start.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 24/7 常駐交易 bot。__PROJECT_DIR__ 由 install.sh 換成實際路徑。 -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingbot.start</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_DIR__/deploy/macos/start_bot.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- 24/7 常駐（取代舊的每日 08:30 StartCalendarInterval）：
         RunAtLoad=登入即啟動；KeepAlive=任何退出（含優雅關閉 code 0）launchd 立即重啟。
         bot 內建 APScheduler(Asia/Taipei) 只在交易日盤中動作（非交易日 is_trading_day() 自動跳過）；
         start_bot.sh 內的 caffeinate 防 idle 睡眠（需保持筆電開蓋，闔蓋仍會強制睡眠）。
         停機維護（如初始建倉）＝launchctl bootout gui/$(id -u)/com.tradingbot.start；
         更新程式碼＝launchctl kickstart -k gui/$(id -u)/com.tradingbot.start。 -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PYTHONUTF8</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/launchd_start.out.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/launchd_start.err.log</string>
</dict>
</plist>
